# scripts/hawk.service.in
[Unit]
Description=HA Web Konsole
Wants=hawk-backend.service
After=network.target hawk-backend.service
PropagatesReloadTo=hawk-backend.service

[Service]
Type=simple
User=hacluster
Group=haclient
WorkingDirectory=@WWW_BASE@/hawk
Restart=always
RestartSec=5
Environment="HAWK_ENV=production"
EnvironmentFile=-/etc/sysconfig/hawk

# Run ExecStartPre with root-permissions
PermissionsStartOnly=true
ExecStartPre=@WWW_BASE@/hawk/bin/generate-ssl-cert
ExecStart=/usr/sbin/hawk-apiserver -key ${HAWK_KEY} -cert ${HAWK_CERT} -port ${HAWK_PORT} -config /etc/hawk/server.json

[Install]
WantedBy=multi-user.target
